# debian/aodh-api.init.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          aodh-api
# Required-Start:    $network $local_fs $remote_fs $syslog
# Required-Stop:     $remote_fs
# Should-Start:      mongodb rabbitmq-server
# Should-Stop:       mongodb rabbitmq-server
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: aodh api
# Description:       aodh api
### END INIT INFO

# Author: Thomas Goirand <zigo@debian.org>
DESC="OpenStack Aodh API"
PROJECT_NAME=aodh
NAME=${PROJECT_NAME}-api
DAEMON=/usr/bin/uwsgi_python
DAEMON_ARGS="uwsgi_python --master --die-on-term --logto /var/log/aodh/aodh-api.log --stats localhost:9314 --http-socket :8042 --wsgi-file /usr/share/aodh-common/app.wsgi"
NO_OPENSTACK_CONFIG_FILE_DAEMON_ARG=yes
NO_OPENSTACK_LOGFILE_DAEMON_ARG=yes
USE_SYSLOG=no
