# debian/aodh-listener.init.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          aodh-listener
# Required-Start:    $network $local_fs $remote_fs $syslog
# Required-Stop:     $remote_fs
# Should-Start:      mongodb rabbitmq-server
# Should-Stop:       mongodb rabbitmq-server
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: aodh listener
# Description:       aodh listner
### END INIT INFO

# Author: Thomas Goirand <zigo@debian.org>
DESC="OpenStack Aodh Listener"
PROJECT_NAME=aodh
NAME=${PROJECT_NAME}-listener
